# files/etc/supervisor/conf.d/grafana.conf
[program:grafana-webapp]
command = grafana-server -config /etc/grafana/grafana.ini -homepath /usr/share/grafana
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
